# host/ebpf-agent/ebpf-agent.service
[Unit]
Description=eBPF Security Monitoring Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ebpf-agent -config /etc/ebpf-agent/config.yaml
Restart=always
RestartSec=10
User=root

# Security hardening
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
